# base/atuin/.config/systemd/user/atuin.service
[Unit]
Description=Atuin daemon

[Service]
Type=simple
ExecStart=/usr/bin/atuin daemon
Restart=always

[Install]
WantedBy=default.target

