# base/env/.config/profile.d/00-sd-env.sh
export $(run-parts /usr/lib/systemd/user-environment-generators | sed '/:$/d; /^$/d' | xargs -d'\n') ENV_LOADED=1

dbus-update-activation-environment --systemd --all

