# base/env/.config/profile.d/01-runtime-wsl.sh
#!/usr/bin/sh

# wsl
if [ "$WSL2" = "1" ]; then
  export WSL_HOST_IP=$(grep -Po '(?<=nameserver ).*' /etc/resolv.conf)
fi

