# base/env/.config/profile.d/05-path.sh
#!/usr/bin/sh

# path

export PATH="$HOME/.dotfiles/bin:$GOPATH/bin:$BUN_INSTALL/bin:$XDG_DATA_HOME/npm/bin:$PNPM_HOME:$CARGO_HOME/bin:$HOME/.local/bin:$HOME/.local/bin/tmuxes:$XDG_DATA_HOME/neovim/bin:$PATH"

