# base/env/.config/profile.d/06-locale.sh
#!/usr/bin/sh

[ -f "$XDG_CONFIG_HOME/locale.conf" ] && . "$XDG_CONFIG_HOME/locale.conf"
