# base/env/.config/profile.d/20-gtk-setting.sh
#!/usr/bin/sh

# gtk
export GTK_MODULES="appmenu-gtk-module:canberra-gtk-module"
export GTK_USE_PORTAL=1
export GDK_USE_XFT=1

