# base/env/.config/profile.d/30-tools-rye.sh
#!/usr/bin/sh

[ -f "$RYE_HOME"/env ] && . "$RYE_HOME"/env

