# base/env/.config/profile.d/30-tools-pager.sh
#!/usr/bin/sh

export LESS="-iMRnsSW -z-4 -x4 --mouse --wheel-lines=3 --incsearch --use-color"
export SYSTEMD_LESS="$LESS"

