# base/env/.config/profile.d/30-tools-cargo.sh
#!/usr/bin/sh

[ -f "$CARGO_HOME/env" ] && . "$CARGO_HOME/env"

