# base/env/.config/profile.d/40-dev-sccache.sh
#!/usr/bin/sh

# sccache
SCCACHE_PATH="`command -v sccache`"

if [ -n "$SCCACHE_PATH" ]; then
  export RUSTC_WRAPPER="$SCCACHE_PATH"
fi

