# base/env/.config/profile.d/40-dev-mise.sh
#!/usr/bin/sh

# Add support for JetBrains IDEs with mise
ln -sf ~/.local/share/mise ~/.asdf

