# base/env/.config/profile.d/50-docker.sh
# fscrypt on systemd-homed can't use docker-rootless
export DOCKER_HOST="unix:///$XDG_RUNTIME_DIR/podman/podman.sock"

